# custom-files/etc/crontabs/root
00 5 * * * sleep 5 && touch /etc/banner && reboot
30 * * * *  /etc/hosts-new/updatehosts.sh   >/dev/null  2>&1 &
